# frontend/deploy/triplesr-frontend.service
[Unit]
Description=Triple SR Travelers Frontend (Next.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/triplesrfrontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
